# templates/supervisor_nginx.conf.j2
[program:nginx]
command=/usr/sbin/nginx
user={{ nginx_user_name }}
autostart=true
autorestart=true
redirect_stderr=true
